# pi-client/pi-wol-client.service
[Unit]
Description=Pi WOL WebSocket Client
After=network.target

[Service]
ExecStart=/home/joe/Development/wol-ws-bridge/pi-client/venv/bin/python /home/joe/Development/wol-ws-bridge/pi-client/pi-client.py
WorkingDirectory=/home/joe/Development/wol-ws-bridge/pi-client
Restart=always
User=joe
Environment=NODE_ENV=production
Environment=RELAY_SERVER=https://wol-ws-bridge.onrender.com

[Install]
WantedBy=multi-user.target
